# osprogram/lib/systemd/system/exec-server.service
[Unit]
Description=excute exec-server

[Service]
ExecStart=/usr/bin/exec-server

[Install]
WantedBy=multi-user.target
