# osprogram/usr/share/os-code/exec-server/exec-server.service
[Unit]
Description=excute exec-server

[Service]
ExecStart=/home/lyg/home_work/exec-server/exec-server

[Install]
WantedBy=multi-user.target
